# ops/supervisor.conf
[program:noelwilson]
directory=/opt/apps/noelwilson/
command=make run
user=root
redirect_stderr=true
stdout_logfile=/var/log/noelwilson.log
